# playbooks/files/kindelia-node.service
[Unit]
Description=Kindelia node

[Service]
Type=simple
ExecStart=/usr/local/bin/kindelia start --testnet --mine

[Install]
WantedBy=default.target
